# project/config/production/crontab/customer_api
* * * * * root /bin/bash /var/www/customer_api/project/tool/production/check_update.sh
